# rpm/SOURCES/etc/init.d/pepProxy
#!/usr/bin/env bash
#
# pepProxy         Start/Stop the PEP Proxy
#
# chkconfig: 2345 99 60
# description: Orion Policy Enforcement Point
### BEGIN INIT INFO
# Provides: pepProxy
# Required-Start: $local_fs $syslog
# Required-Stop: $local_fs $syslog
# Default-Start:  345
# Default-Stop: 90
# Short-Description: run pepProxy
# Description: The Policy Enforcement Point is part of the Access Control system of
# the Fiware Platform. Its main purpose is to filter the access to the Context Broker
# based on policies defined in the Access Control system. The PEP Proxy validates all
# the incoming requests against the Access Control.
### END INIT INFO

export SYSTEMCTL_SKIP_REDIRECT=true

# Copyright 2013 Telefonica Investigacion y Desarrollo, S.A.U
#
# This file is part of Orion Policy Enforcement Point.
#
# Orion Policy Enforcement Point is free software: you can redistribute it and/or
# modify it under the terms of the GNU Affero General Public License as
# published by the Free Software Foundation, either version 3 of the
# License, or (at your option) any later version.
#
# Orion Policy Enforcement Point is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE. See the GNU Affero
# General Public License for more details.
#
# You should have received a copy of the GNU Affero General Public License
# along with Orion Policy Enforcement Point. If not, see http://www.gnu.org/licenses/.
#
# For those usages not covered by this license please contact with
# fermin at tid dot es


. /etc/rc.d/init.d/functions

PARAM=${1}
INSTANCE=${2}
PEP_COMPONENT_NAME="pepProxy"
COMPONENT_PATH="/opt/${PEP_COMPONENT_NAME}"
COMPONENT_EXE="${COMPONENT_PATH}/bin/${PEP_COMPONENT_NAME}"
COMPONENT_CONF_PATH="/etc/${PEP_COMPONENT_NAME}.d"
COMPONENT_PREFIX="pepproxy_"
# The PSINFO stores the output of psinfo function that shown the process area (ps) filtered by a PATTERN regexp
PSINFO=""
# The regexp used for filter in the process area (ps)
PSPATTERN=""
# No path for Node
NODEPATH="node"

# Assure that exist PID_DIR and LOG_DIR
mkdir -p /var/run/${PEP_COMPONENT_NAME} /var/log/${PEP_COMPONENT_NAME}
chown pepproxy.pepproxy /var/run/${PEP_COMPONENT_NAME} /var/log/${PEP_COMPONENT_NAME}


# Function that show the process area (ps) filtered by a regexp (PSPATTERN)
# The first parameter is optional and define the instancename
# Two behavihours
# - Find all results for all instances of one service "${NODEPATH} ${COMPONENT_EXE}"
# - Find all results for one instance of one service "${NODEPATH} ${COMPONENT_EXE} ${1}\$"
function psinfo
{
  [[ "${1}" == "" ]] && PSPATTERN="${NODEPATH} ${COMPONENT_EXE}"
  [[ "${1}" != "" ]] && PSPATTERN="${NODEPATH} ${COMPONENT_EXE} ${1}\$"
  # TODO: For RH7 we will use pgrep -la -f "PATTERN"
  PSINFO="$(pgrep -l -f "^${PSPATTERN}")"
}

component_start()
{

    local result=0
    local run_instance=${1}

    if [[ ! -x ${COMPONENT_EXE} ]]; then
        printf "%s\n" "Fail - missing ${COMPONENT_EXE} executable"
        exit 1
    fi

    if [[ -z ${run_instance} ]]; then
        list_instances="${COMPONENT_CONF_PATH}/${COMPONENT_PREFIX}*.conf"
    else
        list_instances="${COMPONENT_CONF_PATH}/${COMPONENT_PREFIX}${run_instance}.conf"
    fi

    if [[ $(ls -l ${list_instances} 2> /dev/null | wc -l) -eq 0 ]]; then
        if [[ ${run_instance} == "" ]]; then
            printf "%s\n" "There aren't any instance of ${PEP_COMPONENT_NAME} configured. Refer to file ${COMPONENT_CONF_PATH}/README.md for further information."
        else
            printf "%s\n" "There aren't any instance of ${PEP_COMPONENT_NAME} configured with the name ${run_instance}. Refer to file ${COMPONENT_CONF_PATH}/README.md for further information."
        fi
        return 1
    fi

    for _instance in ${list_instances}
    do
    
        local NAME
        NAME=${_instance%.conf}
        NAME=${NAME#*${COMPONENT_PREFIX}}

        set -a
        source ${_instance}

        local LOG_FILE="/var/log/${PEP_COMPONENT_NAME}/${COMPONENT_PREFIX}${NAME}.log"
        local PID_FILE="/var/run/${PEP_COMPONENT_NAME}/${COMPONENT_PREFIX}${NAME}.pid"

        printf "%s" "Starting instance ${NAME} of ${PEP_COMPONENT_NAME}... "

        status -p ${PID_FILE} ${COMPONENT_EXE} &> /dev/null 
        if [[ ${?} -eq 0 ]]; then
            printf "%s\n" " Already running, skipping $(success)"
            continue
        fi

        # Launch one instance with first parameter  being the name of instance
        su $PROXY_USER -p -c "cd ${COMPONENT_PATH}; ${COMPONENT_EXE} ${NAME} &>> ${LOG_FILE} & echo \$! > ${PID_FILE}"
        sleep 2 # some cortesy time to process startup or die
        local PID=$(cat ${PID_FILE})
        local PEP_PID=$(ps -ef | grep -v "grep" | grep "${PID:-not_found}")
        if [[ -z ${PEP_PID} ]]; then
            printf "%s\n" "$(failure)"
            result=$((${result}+1))
            rm -f ${PID_FILE} &> /dev/null 
        else
            chown ${COMPONENT_USER}:${COMPONENT_USER} ${PID_FILE}
            printf "%s\n" "$(success)"
        fi

    done

    return ${result}
    
}

component_stop()
{
    local result=0
    local run_instance=${1}

    if [[ -z ${run_instance} ]]; then
        list_instances="${COMPONENT_CONF_PATH}/${COMPONENT_PREFIX}*.conf"
    else
        list_instances="${COMPONENT_CONF_PATH}/${COMPONENT_PREFIX}${run_instance}.conf"
    fi

    if [[ $(ls -l ${list_instances} 2> /dev/null | wc -l) -eq 0 ]]; then
        printf "%s\n" "There aren't any instance of ${PEP_COMPONENT_NAME} running $(success)"
        return 0
    fi

    for _instance in ${list_instances}
    do

        local NAME
        NAME=${_instance%.conf}
        NAME=${NAME#*${COMPONENT_PREFIX}}

        printf "%s" "Stopping instance ${NAME} of ${PEP_COMPONENT_NAME}... "

        # Obtain process status
        psinfo ${NAME}
        if [ -z "${PSINFO}" ]
        then
          echo "INFO: Instance ${NAME} of ${PEP_COMPONENT_NAME} it is not running. Do not nothing"
        else
          echo "INFO: Stop instance ${NAME} of ${PEP_COMPONENT_NAME}... "
          echo "${PSINFO}"
          echo "${PSINFO}" | awk '{print $1}' | xargs -r kill
          sleep 1
          psinfo ${NAME}
          echo "${PSINFO}" | awk '{print $1}' | xargs -r kill -9
        fi

        # Remove the PID file instance
        rm -f /var/run/${PEP_COMPONENT_NAME}/${COMPONENT_PREFIX}${NAME}.pid &> /dev/null 

    done
    return ${result}
}

component_status()
{
    local result=0
    local run_instance=${1}

    if [[ -z ${run_instance} ]]; then
        list_run_instances="/var/run/${PEP_COMPONENT_NAME}/${COMPONENT_PREFIX}*.pid"
    else
        list_run_instances="/var/run/${PEP_COMPONENT_NAME}/${COMPONENT_PREFIX}${run_instance}.pid"
    fi

    if [[ $(ls -l ${list_run_instances} 2> /dev/null | wc -l) -eq 0 ]]; then
        printf "%s\n" "There aren't any instance of ${PEP_COMPONENT_NAME} running."
        return 1
    fi

    for _instance in ${list_run_instances}
    do

        local NAME
        NAME=${_instance%.pid}
        NAME=${NAME#*${COMPONENT_PREFIX}}

        printf "%s\n" "${PEP_COMPONENT_NAME} instance ${NAME} status... "
        status -p ${_instance} ${NODE_EXEC}
        result=$((${result}+${?}))

    done

    return ${result}
}

case ${PARAM} in

    'start')
        component_start ${INSTANCE}
        ;;

    'stop')
        component_stop ${INSTANCE}
        ;;

    'restart')
        component_stop ${INSTANCE}
        component_start ${INSTANCE}
        ;;

    'status')
        component_status ${INSTANCE}
        ;;

esac

